# tieba-sign/s6-overlay/etc/services.d/crond/run
#!/usr/bin/with-contenv bash

echo "[services.d] cron service running successful."

exec s6-setuidgid root crond -f &>/dev/null
